# default/replication.cron
SHELL=/bin/bash
BASH_ENV=/noninteractive.bash_env
0 * * * * /usr/local/bin/replication.sh
